# deploy/systemd/fortress-dashboard.service
[Unit]
Description=Fortress Guest Platform Dashboard
After=network-online.target fortress-backend.service
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/Fortress-Prime/fortress-guest-platform/frontend-next
EnvironmentFile=-/home/admin/Fortress-Prime/fortress-guest-platform/.env
EnvironmentFile=-/home/admin/Fortress-Prime/fortress-guest-platform/.env.dgx
EnvironmentFile=-/home/admin/Fortress-Prime/.env.security
ExecStart=/usr/local/bin/run-fortress-dashboard.sh
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
