# deploy/systemd/fortress-arq-worker.service.d/00-secrets.conf
# Drop-in: load runtime secrets from `pass` before the worker starts.
#
# Pre-step writes /run/fortress/secrets.env (mode 0600, tmpfs-backed,
# wiped on reboot). EnvironmentFile picks it up so every var lands in
# the worker's process environment with no .env edit required.
#
# Requirement: gpg-agent for the `admin` user must have the GPG key's
# passphrase cached when systemd starts the worker. Configure once via
# ~/.gnupg/gpg-agent.conf:
#     default-cache-ttl 34560000
#     max-cache-ttl     34560000
# then `gpg-connect-agent reloadagent /bye` and prime once with any
# `pass show fortress/...` call as the admin user. See
# deploy/secrets/README.md for the full setup.
#
# A pass failure (missing entry, agent locked, manifest typo) makes
# ExecStartPre fail — systemd refuses to start the main ExecStart and
# logs the variable that failed. The worker never enters the poll loop
# with stale or empty secrets.

[Service]
# Ensure /run/fortress exists with the right perms across reboots.
RuntimeDirectory=fortress
RuntimeDirectoryMode=0700

ExecStartPre=/usr/local/bin/fortress-load-secrets --output /run/fortress/secrets.env

# The `-` prefix makes the file optional during the very first boot
# before the pre-step has run. ExecStartPre still gates startup on
# successful generation.
EnvironmentFile=-/run/fortress/secrets.env
